# com.arjundivecha.neo4j-guard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
    Neo4j stale-pidfile boot guard (2026-06-18).
    Runs scripts/ops/neo4j_guard.sh at login (RunAtLoad) and every 30 minutes.
    If Neo4j's bolt port 7687 is down, it removes any STALE Homebrew PID file
    (the post-hard-crash trap that makes `neo4j console` refuse to start) and
    restarts the Homebrew Neo4j job. No-op when Neo4j is already healthy.
    This is the recurrence fix for the 2026-06-18 outage where a power loss
    left a stale neo4j.pid and the ASADO loop's graph steps all failed.

    INSTALL (the agent does NOT auto-edit ~/Library/LaunchAgents — your call):
      cp "com.arjundivecha.neo4j-guard.plist" ~/Library/LaunchAgents/
      launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.arjundivecha.neo4j-guard.plist
    Logs to ASADO/Data/logs/neo4j_guard.log (script) + neo4j_guard_launchd.log (stdout/err).
    -->
    <key>Label</key>
    <string>com.arjundivecha.neo4j-guard</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/arjundivecha/Dropbox/AAA Backup/A Working/ASADO/scripts/ops/neo4j_guard.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/arjundivecha/Dropbox/AAA Backup/A Working/ASADO</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/arjundivecha/Dropbox/AAA Backup/A Working/ASADO/Data/logs/neo4j_guard_launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/arjundivecha/Dropbox/AAA Backup/A Working/ASADO/Data/logs/neo4j_guard_launchd.log</string>
</dict>
</plist>
